# maintenance/.config/systemd/user/clean-downloads.service
[Unit]
Description=Clean up old downloads
Wants=daily.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/date

[Install]
WantedBy=multi-user.target
